# deploy/ansible/templates/octo.service.j2
# Octo Control Plane Server
# Generated by Ansible

[Unit]
Description=Octo Control Plane Server
After=network.target
{% if octo_mode == "multi" %}
Wants=octo-agent@.service
{% endif %}

[Service]
Type=simple
User={{ octo_user }}
Group={{ octo_group }}

# Working directory
WorkingDirectory=/var/lib/octo

# Environment
Environment=OCTO_CONFIG=/etc/octo/config.toml
Environment=RUST_LOG=info
Environment=PATH=/usr/local/cargo/bin:/usr/local/bin:/usr/bin:/bin

# Create directories
StateDirectory=octo
RuntimeDirectory=octo
RuntimeDirectoryMode=0755
ConfigurationDirectory=octo

# Main process
{% if octo_backend == "local" %}
ExecStart=/usr/local/bin/octo serve --local-mode
{% else %}
ExecStart=/usr/local/bin/octo serve
{% endif %}

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/octo
ReadWritePaths=/run/octo
PrivateTmp=true

# Allow binding to privileged ports (optional, for port 80/443)
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
