# deploy/systemd/octo-user.service
# Octo User Agent - Systemd user service
#
# Alternative to octo-agent@.service for systems using systemd user services.
# This runs as part of the user's session, not as a system service.
#
# Installation (as user):
#   mkdir -p ~/.config/systemd/user
#   cp octo-user.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable octo-user
#   systemctl --user start octo-user
#
# The agent communicates via Unix socket in XDG_RUNTIME_DIR

[Unit]
Description=Octo User Agent
After=default.target

[Service]
Type=simple

# Environment
Environment=OCTO_MODE=user

# Socket path in user's runtime directory
ExecStart=/usr/local/bin/octo-agent serve --socket %t/octo-agent.sock

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
