# deploy/systemd/octo-agent@.service
# Octo Agent - Per-user systemd service
# 
# This service runs the octo-agent for a specific user, providing:
# - OpenCode instance running as that user
# - File server for workspace access
# - Terminal (ttyd) access
#
# Installation (as root):
#   cp octo-agent@.service /etc/systemd/system/
#   systemctl daemon-reload
#
# Usage:
#   systemctl start octo-agent@alice    # Start agent for user alice
#   systemctl enable octo-agent@alice   # Enable on boot
#   systemctl status octo-agent@alice   # Check status
#
# The agent listens on a Unix socket at /run/octo/agent-%i.sock

[Unit]
Description=Octo Agent for %i
After=network.target
PartOf=octo.service

[Service]
Type=simple
User=%i
Group=%i

# Working directory - user's home
WorkingDirectory=~

# Environment
Environment=OCTO_USER=%i
Environment=OCTO_SOCKET=/run/octo/agent-%i.sock
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Create runtime directory
RuntimeDirectory=octo
RuntimeDirectoryMode=0755

# Main process - octo-agent binary
# Listens on Unix socket for commands from octo-server
ExecStart=/usr/local/bin/octo-agent serve --socket /run/octo/agent-%i.sock

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i
ReadWritePaths=/var/lib/octo/users/%i
PrivateTmp=true

[Install]
WantedBy=multi-user.target
