# systemd/octo-runner.socket
# octo-runner socket activation (optional)
#
# This allows the runner to be started on-demand when a connection is made.
# Alternative to running the service constantly.
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp octo-runner.socket ~/.config/systemd/user/
#   cp octo-runner@.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable octo-runner.socket
#   systemctl --user start octo-runner.socket

[Unit]
Description=Octo Runner Socket
Documentation=https://github.com/byteowlz/octo

[Socket]
ListenStream=%t/octo-runner.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
