# deploy/archiso/airootfs/etc/systemd/system/octo-first-boot.service
# Octo First Boot Service
# Runs once on first boot to complete setup
[Unit]
Description=Octo First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/octo/.first-boot-complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/octo-first-boot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
